# weather-display.service
[Unit]
Description=Raspberry Pi Weather Display System
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/skyforge_project/
ExecStart=/usr/bin/python3 /home/pi/skyforge_project/system_controller.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
